# userspace_app/bash/shlug/crontab
0 0     * * 2   root    /dir/to/shlug_cron.sh
